# bsp/arch-am5728/apps/base/etc/inittab
::sysinit:/etc/rcS
::restart:/sbin/init
ttyS2::askfirst:-/bin/sh
